# systemd/user/psd-resync.timer.d/frequency.conf
[Unit]
Description=Timer for profile-sync-daemon

[Timer]
OnUnitActiveSec=
OnUnitActiveSec=10min
